# scripts/com.ptd.aws-sync-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ptd.aws-sync-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/milosvukovic/client-vital-suite/scripts/aws-sync-bridge.cjs</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/ptd-aws-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ptd-aws-sync-error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
